# meta-rdk-matter/recipes-matter/chip-lighting-app/files/chip-lighting-app.service
[Unit]
Description=Chip Lighting App for Matter IoT
After=network-online.target avahi-daemon.service dbus.service ot-br-posix.service
Wants=network-online.target
Requires=avahi-daemon.service

[Service]
Type=simple
ExecStart=/usr/bin/chip-lighting-app
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

