# bluez5/5.66-r0/sysroot-destdir/usr/lib/systemd/user/obex.service
[Unit]
Description=Bluetooth OBEX service

[Service]
Type=dbus
BusName=org.bluez.obex
ExecStart=/usr/lib/bluez5/bluetooth/obexd

[Install]
Alias=dbus-org.bluez.obex.service
